# radlog.service
[Unit]
Description=Geiger Counter Data Logger

[Service]
Type=simple
ExecStart=/opt/radlog/radlog.py
Restart=on-failure
