# radlog.crontab
*/5 * * * * mdwyer /opt/radlog/radpic.sh
